# systemd/user/echoflow.service.in
# SPDX-FileCopyrightText: 2026 Hualet Wang
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=EchoFlow offline voice input service
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/echoflow
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/echoflow-service
Restart=on-failure

[Install]
WantedBy=default.target
