# systemd/user/echoflow-ui.service.in
# SPDX-FileCopyrightText: 2026 Hualet Wang
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=EchoFlow tooltip UI host
After=graphical-session.target

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/echoflow-ui
Restart=on-failure

[Install]
WantedBy=default.target
